# data/synthetic/cbse_attack/attack_logs/cbse_2026_breach.json
[
  {
    "event_id": "MAIL_DELIVERY",
    "sequence": 1,
    "timestamp": "2026-03-02T11:33:14+05:30",
    "host": "CBSE-MAIL-GW-01",
    "target_account": "exam_admin_pankaj",
    "event_type": "phishing_email_delivery",
    "email_from": "noreply@cbse-results.in",
    "email_subject": "URGENT: Board Exam 2026 Final Paper Upload — Action Required Before 5 PM",
    "email_body_preview": "Dear Controller of Examinations, This is a mandatory action required for Class 12 paper upload validation. Your credentials must be re-verified due to security audit. Click link below to confirm access...",
    "malicious_link": "https://cbse-results.in/admin/verify?token=aGFja2VkYnlhcHQ0MQ==",
    "link_resolves_to": "45.142.212.100",
    "attachment": null,
    "mail_security_verdict": "PASSED_SPAM_FILTER",
    "spf_check": "PASS",
    "dkim_check": "PASS",
    "sender_domain_age_days": 8,
    "mitre_technique": "T1566.002",
    "mitre_technique_name": "Phishing: Spearphishing Link",
    "anomaly_signals": [
      "lookalike_domain: cbse-results.in vs legitimate cbse.gov.in — subtle .in vs .gov.in substitution",
      "domain_age: registered 8 days ago — brand new domain impersonating government authority",
      "credential_harvesting_url: /admin/verify pattern with base64 token — credential phishing page",
      "urgency_social_engineering: 'Before 5 PM' + 'URGENT' + 'mandatory' — pressure tactics for credential bypass",
      "exam_period_timing: 3 weeks before board exams — maximum pressure on exam administration staff"
    ],
    "ground_truth_note": "MALICIOUS: Spearphishing email targeting Pankaj Kumar, CBSE Exam Paper Administrator. The lookalike domain 'cbse-results.in' passes SPF/DKIM checks because the attacker actually configured valid DNS records for their phishing domain — a sign of sophisticated preparation. The base64 URL token decodes to 'hackedbyapt41'. Email was sent 3 weeks before Class 12 board exams — maximum pressure period when exam staff are most likely to act quickly on 'urgent' requests. The link leads to a pixel-perfect clone of the CBSE admin portal at a server in the Netherlands (ASN: FRANTECH SOLUTIONS)."
  },
  {
    "event_id": "4624",
    "sequence": 2,
    "timestamp": "2026-03-02T11:47:22+05:30",
    "host": "CBSE-EXAM-PORTAL-01",
    "account": "exam_admin_pankaj",
    "logon_type": 10,
    "source_ip": "45.142.212.100",
    "dest_ip": "10.1.5.20",
    "auth_package": "Forms_Authentication",
    "workstation_name": "PHISHING_PROXY",
    "session_id": "sess_hijacked_4f2a",
    "mitre_technique": "T1078",
    "mitre_technique_name": "Valid Accounts",
    "anomaly_signals": [
      "foreign_ip_logon: 45.142.212.100 is Netherlands IP — CBSE admin portal never accessed from outside India",
      "real_time_credential_relay: attacker relays credentials in real-time from phishing page to actual portal",
      "source_ip_never_seen: zero baseline entries for any 45.x.x.x source IP accessing exam portal",
      "session_immediately_suspicious: session from Netherlands on government exam portal = immediate red flag"
    ],
    "ground_truth_note": "MALICIOUS: Real-time credential relay attack (adversary-in-the-middle phishing). When exam_admin_pankaj entered his credentials on the fake CBSE portal, the attacker's server immediately relayed those credentials to the real CBSE portal, obtaining a valid session. This is an evilginx2-style attack — the phishing page acts as a transparent proxy. The resulting session appears legitimate to the server but originates from 45.142.212.100 (Netherlands) — a geographic impossibility for a CBSE employee working in Delhi."
  },
  {
    "event_id": "S3_BULK_ACCESS",
    "sequence": 3,
    "timestamp": "2026-03-02T11:52:08+05:30",
    "host": "CBSE-EXAM-PORTAL-01",
    "account": "exam_admin_pankaj",
    "session_id": "sess_hijacked_4f2a",
    "event_type": "bulk_object_download",
    "source_ip": "45.142.212.100",
    "storage_system": "CBSE Secure Exam Paper Repository",
    "storage_path": "/exam-papers/2026/class-12/",
    "objects_accessed": [
      {"filename": "class12_mathematics_set1_2026.pdf", "size_mb": 4.2, "classification": "TOP SECRET — EXAM PAPER"},
      {"filename": "class12_mathematics_set2_2026.pdf", "size_mb": 4.1, "classification": "TOP SECRET — EXAM PAPER"},
      {"filename": "class12_mathematics_set3_2026.pdf", "size_mb": 4.3, "classification": "TOP SECRET — EXAM PAPER"},
      {"filename": "class12_physics_set1_2026.pdf", "size_mb": 3.8, "classification": "TOP SECRET — EXAM PAPER"},
      {"filename": "class12_physics_set2_2026.pdf", "size_mb": 3.9, "classification": "TOP SECRET — EXAM PAPER"},
      {"filename": "class12_chemistry_set1_2026.pdf", "size_mb": 4.0, "classification": "TOP SECRET — EXAM PAPER"},
      {"filename": "class12_chemistry_set2_2026.pdf", "size_mb": 3.7, "classification": "TOP SECRET — EXAM PAPER"},
      {"filename": "class12_biology_set1_2026.pdf", "size_mb": 4.5, "classification": "TOP SECRET — EXAM PAPER"},
      {"filename": "class12_english_core_set1_2026.pdf", "size_mb": 2.8, "classification": "TOP SECRET — EXAM PAPER"},
      {"filename": "class12_hindi_core_set1_2026.pdf", "size_mb": 3.1, "classification": "TOP SECRET — EXAM PAPER"}
    ],
    "total_files_downloaded": 10,
    "total_size_mb": 38.4,
    "download_duration_seconds": 47,
    "mitre_technique": "T1530",
    "mitre_technique_name": "Data from Cloud Storage Object",
    "anomaly_signals": [
      "bulk_download_10_papers: 10 exam papers downloaded in 47 seconds — no admin downloads all papers simultaneously",
      "foreign_ip_download: exam papers downloaded to 45.142.212.100 (Netherlands) — never before seen in access logs",
      "top_secret_data: class 12 board exam papers are national security equivalent — 14 million students taking exams",
      "access_time: 5 minutes after the phishing credential theft — same session, immediate exploitation",
      "no_business_justification: no workflow requires downloading all papers at once to an external IP"
    ],
    "ground_truth_note": "MALICIOUS: Within 5 minutes of credential theft, the attacker downloads 10 exam papers for the 2026 Class 12 board exams — affecting 14 million students. This is the primary objective: selling leaked exam papers for financial gain or national disruption. The 47-second bulk download indicates automated tooling. The attacker now has ALL major subjects: Math, Physics, Chemistry, Biology, English, Hindi. If published before the exam, this would require reprinting and rescheduling all Class 12 exams (affecting 14 million students, massive logistics challenge for CBSE)."
  },
  {
    "event_id": "4688",
    "sequence": 4,
    "timestamp": "2026-03-14T02:00:00+05:30",
    "host": "CBSE-STUDENT-DB-01",
    "account": "SYSTEM",
    "process_name": "svhost64.exe",
    "command_line": "svhost64.exe --encrypt --target D:\\StudentRegistration\\ --extension .CBSE_LOCKED --ransom-note READ_STUDENTS.txt",
    "process_parent": "WinUpdateHelper64",
    "mitre_technique": "T1486",
    "mitre_technique_name": "Data Encrypted for Impact",
    "affected_database": "CBSE Student Registration Database 2026",
    "records_targeted": 14200000,
    "timing": "12 days after initial access — night before Class 12 board exams begin",
    "anomaly_signals": [
      "ransomware_day_before_exam: encryption triggered night before board exams — maximum disruption leverage",
      "14_million_records: student registration database for 14 million Class 12 students",
      "svhost64_masquerade: 'svhost64.exe' mimics legitimate svchost.exe — typo masquerade",
      "WinUpdateHelper64_parent: persistence established via scheduled task during dwell period",
      "exam_timing_precision: 02:00 AM start ensures disruption discovered exactly at school opening time"
    ],
    "ground_truth_note": "MALICIOUS: Ransomware deployed on the student registration database the night before Class 12 board exams begin. 14.2 million student records encrypted — this means CBSE cannot verify student identities, admit cards, or exam center assignments for any student. The timing (02:00 AM, night before exam day 1) is deliberately chosen: discovered at 06:00-07:00 AM when exam centers open, leaving zero time for recovery before students arrive. The double extortion: (1) exam papers already stolen and threatening to sell/publish, (2) student registration database encrypted. Maximum national disruption achieved."
  },
  {
    "event_id": "4698",
    "sequence": 5,
    "timestamp": "2026-03-02T12:15:33+05:30",
    "host": "CBSE-EXAM-PORTAL-01",
    "account": "exam_admin_pankaj",
    "session_id": "sess_hijacked_4f2a",
    "event_type": "scheduled_task_created",
    "task_name": "WinUpdateHelper64",
    "task_path": "\\Microsoft\\Windows\\Application Experience\\",
    "task_command": "C:\\ProgramData\\Microsoft\\Windows\\WinUpdate64\\svhost64.exe",
    "task_trigger": "OneTime",
    "task_trigger_time": "2026-03-14T02:00:00+05:30",
    "task_run_as": "SYSTEM",
    "mitre_technique": "T1053.005",
    "anomaly_signals": [
      "12_day_delayed_trigger: task scheduled to run 12 days in future — time-bomb persistence mechanism",
      "exam_date_target: 2026-03-14 is first day of Class 12 board exams — precisely timed",
      "malicious_executable: svhost64.exe in ProgramData (not System32) — binary staged after credential theft",
      "admin_creating_system_task: hijacked admin session creating SYSTEM-level task from Netherlands IP"
    ],
    "ground_truth_note": "MALICIOUS: A time-bomb scheduled task created during the hijacked session, set to trigger exactly on the first day of Class 12 board exams. This reveals extraordinary pre-planning: the attacker knew the CBSE exam schedule, chose the precise maximum-impact date, and created a delayed-execution mechanism to fire 12 days after the initial access. The task is created from the hijacked Netherlands IP session — an impossible action for legitimate CBSE IT operations. Detection of this task creation would have given 12 days to remove the time-bomb and prevent exam day disruption."
  },
  {
    "event_id": "DNS_QUERY",
    "sequence": 6,
    "timestamp": "2026-03-02T12:03:44+05:30",
    "host": "CBSE-EXAM-PORTAL-01",
    "source_ip": "10.1.5.20",
    "query_domain": "cdn.microsott-update.com",
    "query_type": "A",
    "response_ip": "45.142.212.100",
    "response_code": "NOERROR",
    "ttl": 60,
    "mitre_technique": "T1071.001",
    "anomaly_signals": [
      "typosquat: 'microsott-update.com' (double t, hyphen) vs 'microsoft.com'",
      "same_tor_operator: 45.142.212.100 in same ASN as AIIMS C2 (FRANTECH SOLUTIONS)",
      "low_ttl: TTL 60 seconds — fast-flux C2 infrastructure",
      "cbse_portal_external_dns: government exam portal resolving non-government domains"
    ],
    "ground_truth_note": "MALICIOUS: C2 domain for CBSE attack using same infrastructure pattern as AIIMS attack (typosquat + FRANTECH SOLUTIONS ASN + 60s TTL). This is a critical forensic link between the two incidents — same threat actor, same infrastructure playbook. The CBSE exam portal should never resolve external Microsoft-lookalike domains. This DNS query, cross-referenced against the AIIMS incident IOCs, would have identified the same attacker group within minutes."
  },
  {
    "event_id": "4663",
    "sequence": 7,
    "timestamp": "2026-03-02T12:20:11+05:30",
    "host": "CBSE-EXAM-PORTAL-01",
    "account": "exam_admin_pankaj",
    "session_id": "sess_hijacked_4f2a",
    "event_type": "binary_upload",
    "source_ip": "45.142.212.100",
    "file_written": "C:\\ProgramData\\Microsoft\\Windows\\WinUpdate64\\svhost64.exe",
    "file_size_bytes": 892416,
    "file_hash_sha256": "a94f5f3a2d8e4b1c9f7e6d2a1b3c5e8f9a0b2d4c6e8f0a2b4d6e8f0a2b4d6e8",
    "mitre_technique": "T1105",
    "mitre_technique_name": "Ingress Tool Transfer",
    "anomaly_signals": [
      "binary_upload_from_netherlands: executable uploaded from 45.142.212.100 via stolen admin session",
      "non_standard_path: ProgramData\\Microsoft\\Windows\\WinUpdate64\\ — directory doesn't exist by default, created by attacker",
      "pe_executable: svhost64.exe is PE32+ executable — uploaded malware binary",
      "upload_via_web_session: file uploaded through web application session — application layer ingress"
    ],
    "ground_truth_note": "MALICIOUS: The ransomware binary (svhost64.exe, 872KB PE32+ executable) is uploaded directly through the hijacked exam admin web session. The attacker creates a fake Windows directory 'WinUpdate64' in ProgramData to store the binary. The SHA256 hash matches known FinFisher/APT41-variant ransomware payloads. Upload via a legitimate web session bypasses many network controls — the upload appears as legitimate web traffic from the admin portal."
  },
  {
    "event_id": "4688",
    "sequence": 8,
    "timestamp": "2026-03-14T02:45:00+05:30",
    "host": "CBSE-STUDENT-DB-01",
    "account": "SYSTEM",
    "process_name": "notepad.exe",
    "command_line": "notepad.exe D:\\StudentRegistration\\READ_STUDENTS.txt",
    "event_type": "ransom_note_display",
    "ransom_note_content": "CBSE BOARD EXAM 2026 — SYSTEM COMPROMISED\n\nAll student registration records for Class 12 Board Exams 2026 have been encrypted.\n14,200,000 student records are inaccessible.\n\nWe also possess copies of ALL Class 12 exam papers (Mathematics, Physics, Chemistry, Biology, English, Hindi).\n\nYour options:\n1. PAY RANSOM: USD 8,000,000 in Bitcoin within 48 hours — receive decryption key + papers destroyed\n2. DO NOT PAY: Papers published on Telegram within 48 hours + student data sold to dark web markets\n\nContact: cbse_crisis@tutanota.com\nBitcoin: bc1qm9yr4m9m4r9qv8w4y6v9w2y6m9r4q8v2m\n\nThis is not a drill. Exams begin in 4 hours.",
    "mitre_technique": "T1486",
    "anomaly_signals": [
      "exam_day_timing: ransom note displayed at 03:45 AM, exams begin at 08:00 AM — 4 hour window for response",
      "dual_extortion: both data encryption AND exam paper leak threat — double leverage",
      "14_million_students: impact scale — Class 12 board exams affect entire school-leaving cohort",
      "government_crisis: national education infrastructure disruption — potential for national exam postponement"
    ],
    "ground_truth_note": "MALICIOUS: The ransom note surfaces at 03:45 AM on exam day — 4.25 hours before exams begin. The dual extortion (encrypted student records + threat to leak exam papers already stolen 12 days prior) creates maximum pressure. CBSE officials discovered the encrypted databases at 05:30 AM when IT staff arrived for exam day preparation. With 4 hours to exam start and 14.2 million students arriving at exam centers, the government faced an impossible choice: delay all exams (national crisis) or proceed without functional student verification systems. This demonstrates that detecting the initial phishing credential theft (Day 1) or the scheduled task creation (Day 1, 30 minutes after phishing) would have prevented the entire exam disruption."
  }
]
